# opt/logcenter/share/systemd/rsyslog.service.d/logcenter.conf
[Service]
ExecStart=
ExecStart=/opt/logcenter/bin/logcenter-env /usr/sbin/rsyslogd -n $SYSLOGD_OPTIONS
