# opt/logcenter/share/systemd/metricbeat.service.d/logcenter.conf
[Service]
Environment=BEAT_LOG_OPTS="--environment default"
