# opt/logcenter/share/systemd/elasticsearch.service.d/logcenter.conf
[Unit]
StartLimitIntervalSec=0

[Service]
EnvironmentFile=-/etc/sysconfig/elasticsearch.local
LimitMEMLOCK=infinity
Restart=on-failure
RestartSec=10
