# opt/logcenter/share/systemd/kibana.service.d/logcenter.conf
[Service]
Environment=HOSTNAME="%H"
#AmbientCapabilities=CAP_NET_BIND_SERVICE
