# etc/systemd/system/rsyslog.service.d/override.conf
[Service]
EnvironmentFile=-/etc/logcenter/logcenter.conf
